# deploy/systemd/trend-kr.timer
[Unit]
Description=Daily KR pipeline at 20:30 KST

[Timer]
OnCalendar=*-*-* 20:30:00 Asia/Seoul
Persistent=true
Unit=trend-kr.service

[Install]
WantedBy=timers.target
