# deploy/systemd/trend-us.timer
[Unit]
Description=Daily US pipeline at 09:30 KST

[Timer]
OnCalendar=*-*-* 09:30:00 Asia/Seoul
Persistent=true
Unit=trend-us.service

[Install]
WantedBy=timers.target
