# deploy/systemd/trend-kr.timer
[Unit]
Description=Daily KR pipeline at 08:30 KST

[Timer]
OnCalendar=*-*-* 08:30:00 Asia/Seoul
Persistent=true
Unit=trend-kr.service

[Install]
WantedBy=timers.target
