# deploy/systemd/trend-us.service
[Unit]
Description=Trend Trading US Pipeline
After=network-online.target

[Service]
Type=oneshot
User=trendtrading
Group=trendtrading
WorkingDirectory=/home/trendtrading/Trend_Trading
EnvironmentFile=/home/trendtrading/Trend_Trading/config/.env
Environment=PYTHONPATH=/home/trendtrading/Trend_Trading
ExecStart=/home/trendtrading/Trend_Trading/venv/bin/python3 update_all.py --market us
StandardOutput=append:/home/trendtrading/Trend_Trading/logs/daily_us.log
StandardError=append:/home/trendtrading/Trend_Trading/logs/daily_us.log
